# etc/demozoo-productionvm-supervisor.conf
[program:demozoo]
command=/home/demozoo/.virtualenvs/demozoo/bin/uwsgi --ini etc/demozoo-uwsgi.ini
directory=/home/demozoo/demozoo
autostart=true
autorestart=true
redirect_stderr=True
user=demozoo
