# etc/zxdemo-productionvm-supervisor.conf
[program:zxdemo]
command=/home/demozoo/.virtualenvs/demozoo/bin/uwsgi --ini etc/zxdemo-uwsgi.ini
directory=/home/demozoo/demozoo
autostart=true
autorestart=true
redirect_stderr=True
user=demozoo
